# signal_system/services/macos/com.signalsystem.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.signalsystem.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/YOUR_USERNAME/signal_system/python/listener.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/signal_system/python</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/signal_system/python</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/signal_system/logs/listener.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/signal_system/logs/listener.error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Wait for bridge to start first -->
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

</dict>
</plist>
