# services/macos/com.signalsystem.tv-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- ── Identity ────────────────────────────────────────────── -->
    <key>Label</key>
    <string>com.signalsystem.tv-mcp</string>

    <!-- ── Command ─────────────────────────────────────────────── -->
    <!-- Runs the F3 Streamable HTTP MCP daemon from the local fork. -->
    <!-- __NODE_BIN__ resolves to the absolute node binary path at    -->
    <!-- render time (launchd ProgramArguments[0] needs absolute path).-->
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>/Users/YOUR_USERNAME/tradingview-mcp-aurum/src/server.js</string>
    </array>

    <!-- ── Working directory ────────────────────────────────────── -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/tradingview-mcp-aurum</string>

    <!-- ── Environment ─────────────────────────────────────────── -->
    <!-- PYTHONUNBUFFERED is the anchor inject_env_vars() looks for; -->
    <!-- harmless for node. .env values (including MCP_TRANSPORT,    -->
    <!-- MCP_HTTP_PORT, MCP_TRACE_FILE) are injected after this key. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- ── Logging ──────────────────────────────────────────────── -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/signal_system/logs/tv-mcp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/signal_system/logs/tv-mcp.error.log</string>

    <!-- ── Auto-start on login ──────────────────────────────────── -->
    <key>RunAtLoad</key>
    <true/>

    <!-- ── Auto-restart on crash ───────────────────────────────── -->
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- ── Restart throttle (wait 10s before restarting) ──────── -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

</dict>
</plist>
