# signal_system/services/linux/signal-athena.service
[Unit]
Description=Signal System — ATHENA Dashboard API
After=network.target signal-bridge.service
Requires=signal-bridge.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/signal_system/python
EnvironmentFile=/home/YOUR_USERNAME/signal_system/.env
Environment=PYTHONPATH=/home/YOUR_USERNAME/signal_system/python
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/signal_system/python/athena_api.py
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signal-athena

[Install]
WantedBy=multi-user.target
