# signal_system/services/linux/signal-listener.service
# ── signal-listener.service ───────────────────────────────────────
[Unit]
Description=Signal System — LISTENER Telegram Signal Reader
After=network.target signal-bridge.service
Requires=signal-bridge.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/signal_system/python
EnvironmentFile=/home/YOUR_USERNAME/signal_system/.env
Environment=PYTHONPATH=/home/YOUR_USERNAME/signal_system/python
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/signal_system/python/listener.py
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signal-listener

[Install]
WantedBy=multi-user.target
